# net-misc/eg25-manager/files/gpsd-pinephone.path
[Unit]
Description=Pinephone modem path

[Path]
PathExists=/dev/dummy
PathExists=/dev/ttyUSB1

[Install]
WantedBy=multi-user.target
