# app-mobilephone/pinephone-modem-scripts/files/gpsd-pinephone.path
[Unit]
Description=Pinephone modem path

[Path]
PathExists=/dev/dummy
PathExists=/dev/EG25.NMEA

[Install]
WantedBy=multi-user.target
