# app-mobilephone/pinephone-modem-scripts/files/pinephone-modem-scripts.pinephone-modem-setup.path
[Unit]
Description=Pinephone modem path

[Path]
PathExists=/dev/EG25.AT

[Install]
WantedBy=multi-user.target
