# app-mobilephone/pinephone-modem-scripts/files/pinephone-modem-scripts.pinephone-modem.path
[Unit]
Description=Pinephone modem path

[Path]
PathExists=/dev/dummy
PathExists=/sys/class/modem-power/modem-power/device/powered

[Install]
WantedBy=multi-user.target
